# ipv4-api.service
[Unit]
Description=IP API Service (IPv4)
After=network.target
Documentation=https://github.com/xFanexx/ip-api

[Service]
Type=simple
User=www-data
Group=www-data
ExecStart=/usr/local/bin/ip-api --port 7111
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target